# deploy/systemd/pc28touzhu-source-sync.service
[Unit]
Description=pc28touzhu Source Sync Worker
After=network.target pc28touzhu-platform.service

[Service]
Type=simple
WorkingDirectory=/www/wwwroot/pc28touzhu
EnvironmentFile=/www/wwwroot/pc28touzhu/.env
ExecStart=/usr/bin/env SOURCE_SYNC_ONCE=false /www/wwwroot/pc28touzhu/deploy/systemd/pc28-python.sh platform_source_sync.py
Restart=always
RestartSec=5
User=www
Group=www

[Install]
WantedBy=multi-user.target
