# deploy/systemd/pc28touzhu-platform.service
[Unit]
Description=pc28touzhu Platform API
After=network.target

[Service]
Type=simple
WorkingDirectory=/www/wwwroot/pc28touzhu
EnvironmentFile=/www/wwwroot/pc28touzhu/.env
ExecStart=/usr/bin/env /www/wwwroot/pc28touzhu/deploy/systemd/pc28-python.sh -m pc28touzhu.main
Restart=on-failure
RestartSec=2
User=www
Group=www

[Install]
WantedBy=multi-user.target
