# deploy/systemd/pc28touzhu-telegram-executor@.service
[Unit]
Description=pc28touzhu Telegram Executor (%i)
After=network.target pc28touzhu-platform.service

[Service]
Type=simple
WorkingDirectory=/www/wwwroot/pc28touzhu
Environment=PYTHONPATH=src
EnvironmentFile=/www/wwwroot/pc28touzhu/.env
ExecStart=/usr/bin/env EXECUTOR_ID=%i ONCE=false python3 telegram_executor.py
Restart=always
RestartSec=2
User=www
Group=www

[Install]
WantedBy=multi-user.target
