# deploy/systemd/pc28touzhu-pc28-auto-settlement.service
[Unit]
Description=pc28touzhu PC28 Auto Settlement Worker
After=network.target pc28touzhu-platform.service

[Service]
Type=simple
WorkingDirectory=/www/wwwroot/pc28touzhu
EnvironmentFile=/www/wwwroot/pc28touzhu/.env
Environment=PC28_AUTO_SETTLEMENT_ENABLED=true
ExecStart=/usr/bin/env PC28_AUTO_SETTLEMENT_ONCE=false /www/wwwroot/pc28touzhu/deploy/systemd/pc28-python.sh -m pc28touzhu.pc28_auto_settlement_worker
Restart=always
RestartSec=5
User=www
Group=www

[Install]
WantedBy=multi-user.target
