# deploy/systemd/pc28touzhu-alert-notifier.service
[Unit]
Description=pc28touzhu Platform Alert Notifier (Telegram Bot)
After=network.target pc28touzhu-platform.service

[Service]
Type=simple
WorkingDirectory=/www/wwwroot/pc28touzhu
EnvironmentFile=/www/wwwroot/pc28touzhu/.env
Environment=ALERT_TELEGRAM_ENABLED=true
ExecStart=/usr/bin/env ALERT_NOTIFIER_ONCE=false /www/wwwroot/pc28touzhu/deploy/systemd/pc28-python.sh platform_alert_notifier.py
Restart=always
RestartSec=5
User=www
Group=www

[Install]
WantedBy=multi-user.target
